# systemd/tailscaled.service.d/10-mb-pst-state.conf
# MagicBridge drop-in for the (upstream) tailscaled.service.
#
# PiKVM's rootfs is read-only, so tailscaled's default state dir
# (/var/lib/tailscale) can't be written — the daemon fails with "State store
# failed to initialize ... read-only file system" and Tailscale never works.
#
# Fix: keep tailscaled's state on the PST (Persistent STorage) partition, which
# is writable, and bind-mount it over the stock path so no ExecStart override is
# needed. State then survives reboots (the node stays authenticated).
#
# Harmless if tailscaled isn't installed yet — systemd ignores drop-ins for
# absent units, and it takes effect the moment `tailscale install` creates the unit.
[Unit]
RequiresMountsFor=/var/lib/kvmd/pst

[Service]
ExecStartPre=-/bin/mount -o remount,rw /var/lib/kvmd/pst
ExecStartPre=/bin/sh -c 'mkdir -p /var/lib/kvmd/pst/data/tailscale && chmod 700 /var/lib/kvmd/pst/data/tailscale'
ExecStartPre=/bin/sh -c 'mountpoint -q /var/lib/tailscale || mount --bind /var/lib/kvmd/pst/data/tailscale /var/lib/tailscale'
